# deploy/gcp/duckdns-update.timer
[Unit]
Description=Hourly DuckDNS refresh
Requires=duckdns-update.service

[Timer]
OnBootSec=2min
OnUnitActiveSec=1h
RandomizedDelaySec=5min
Persistent=true

[Install]
WantedBy=timers.target
