# deploy/gcp/duckdns-update.service
[Unit]
Description=Refresh DuckDNS A record for license server
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/duckdns-update.sh
