# deploy/gcp/yg-license-retry-webhooks.service
[Unit]
Description=YgLicenseServer webhook retry queue drain
After=network-online.target yg-license-server.service
Wants=network-online.target

[Service]
Type=oneshot
# --entrypoint python overrides the image's docker-entrypoint.sh (which would
# otherwise prepend `uvicorn app.main:app ...` and trip on the `-m` flag).
# Side benefit: skips the alembic migration step at container start, which
# the main service has already run.
ExecStart=/usr/bin/docker run --rm \
  --env-file /etc/yg-license-server/yg-license-server.env \
  -v /var/lib/yg-license-server:/data \
  --entrypoint python \
  ghcr.io/why-gee/yg-license-server:latest \
  -m app.scripts.retry_webhooks
